# itamae/cookbooks/pacman-syuw/files/etc/systemd/system/pacman-syuw.timer
[Unit]
Description=Download pacman packages daily

[Timer]
OnCalendar=*-*-* 12:17
Persistent=true

[Install]
WantedBy=timers.target
