# itamae/cookbooks/pacman-syuw/files/etc/systemd/system/pacman-syuw.service
[Unit]
Description=Download pacman packages

[Service]
Type=oneshot
User=root
Group=root
ExecStart=pacman -Syuw --noconfirm
