# itamae/cookbooks/ruby-trunk-build/files/etc/systemd/system/ruby-trunk-build.service
[Unit]
Description=Build Ruby trunk

[Service]
Type=oneshot
User=pocke
Group=users
ExecStart=/opt/ruby-trunk-build/ruby-trunk-build.sh
