# itamae/cookbooks/ruby-trunk-build/files/etc/systemd/system/ruby-trunk-build.timer
[Unit]
Description=Build Ruby trunk daily

[Timer]
OnCalendar=*-*-* 04:13
Persistent=true

[Install]
WantedBy=timers.target
