# itamae/cookbooks/ruby-trunk-build/files/Library/LaunchDaemons/launched.ruby-trunk-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>launched.ruby-trunk-build</string>
    <key>ProgramArguments</key>
    <array>
      <string>bash</string>
      <string>/opt/ruby-trunk-build/ruby-trunk-build.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>13</integer>
      </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/kuwabara.masataka/.ruby-trunk-build.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/kuwabara.masataka/.ruby-trunk-build.err</string>
    <key>UserName</key>
    <string>kuwabara.masataka</string>
    <key>ExitTimeOut</key>
    <integer>1200</integer>
    <key>SoftResourceLimits</key>
    <dict>
        <key>CPU</key>
        <integer>1200</integer>
        <key>FileSize</key>
        <integer>2000000000</integer>
        <key>Data</key>
        <integer>2000000000</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>CPU</key>
        <integer>1200</integer>
        <key>FileSize</key>
        <integer>2000000000</integer>
        <key>Data</key>
        <integer>2000000000</integer>
    </dict>
  </dict>
</plist>
